# detections/osquery/macos_persistence_pack.conf
{
  "platform": "darwin",
  "version": "1.0.0",
  "queries": {
    "non_apple_launch_agents": {
      "query": "SELECT la.name, la.label, la.program, la.path, la.run_at_load, la.keep_alive, la.username, h.sha256 FROM launchd la LEFT JOIN hash h ON la.program = h.path WHERE la.label NOT LIKE 'com.apple.%' AND la.path NOT LIKE '/System/%';",
      "interval": 3600,
      "description": "T1543.001/T1543.004: Non-Apple LaunchAgents and LaunchDaemons",
      "value": "Identifies persistence mechanisms from third-party or malicious software"
    },
    "suid_binaries": {
      "query": "SELECT path, username, groupname, permissions, sha256 FROM suid_bin WHERE path NOT LIKE '/usr/bin/%' AND path NOT LIKE '/usr/sbin/%' AND path NOT LIKE '/System/%';",
      "interval": 3600,
      "description": "T1548.001: Non-standard SUID binaries outside system paths",
      "value": "Identifies potential privilege escalation vectors"
    },
    "login_items": {
      "query": "SELECT name, path, type FROM startup_items WHERE path NOT LIKE '/System/%' AND path NOT LIKE '/Applications/%';",
      "interval": 3600,
      "description": "T1547.015: Non-standard login/startup items",
      "value": "Detects persistence via login items"
    },
    "unsigned_running_binaries": {
      "query": "SELECT p.pid, p.name, p.path, p.uid, s.authority, s.identifier FROM processes p LEFT JOIN signature s ON p.path = s.path WHERE s.authority IS NULL OR s.authority = '' LIMIT 50;",
      "interval": 1800,
      "description": "T1553.002: Running processes without valid code signatures",
      "value": "Identifies potentially malicious unsigned processes"
    },
    "dyld_injection_detection": {
      "query": "SELECT pid, name, key, value FROM process_envs WHERE key = 'DYLD_INSERT_LIBRARIES' OR key = 'DYLD_LIBRARY_PATH' OR key LIKE 'DYLD_PRINT_%';",
      "interval": 600,
      "description": "T1574.006: Processes with DYLD injection environment variables",
      "value": "Detects runtime library injection via environment variables"
    },
    "listening_ports_non_apple": {
      "query": "SELECT lp.pid, lp.port, lp.protocol, lp.address, p.name, p.path FROM listening_ports lp JOIN processes p ON lp.pid = p.pid WHERE p.path NOT LIKE '/usr/%' AND p.path NOT LIKE '/System/%' AND lp.port < 65535;",
      "interval": 600,
      "description": "Discovery: Non-system processes with open listening ports",
      "value": "Identifies potential backdoors or C2 listeners"
    },
    "recently_modified_launch_items": {
      "query": "SELECT path, filename, mtime, atime FROM file WHERE (directory = '/Library/LaunchAgents/' OR directory = '/Library/LaunchDaemons/' OR directory LIKE '%/LaunchAgents/') AND mtime > (strftime('%s','now') - 86400);",
      "interval": 3600,
      "description": "T1543: Recently modified LaunchAgent/Daemon plists (last 24h)",
      "value": "Detects newly installed persistence mechanisms"
    },
    "browser_extensions": {
      "query": "SELECT uid, name, path, identifier, version FROM chrome_extensions WHERE path NOT LIKE '%Google%' UNION SELECT uid, name, path, identifier, version FROM safari_extensions WHERE path NOT LIKE '%Apple%' LIMIT 50;",
      "interval": 86400,
      "description": "T1176: Non-standard browser extensions",
      "value": "Identifies potentially malicious browser extensions"
    },
    "tcc_full_disk_access": {
      "query": "SELECT service, client, auth_value, last_modified FROM access WHERE service = 'kTCCServiceSystemPolicyAllFiles' AND auth_value = 2;",
      "interval": 3600,
      "platform": "darwin",
      "description": "T1548: Applications with Full Disk Access TCC permission",
      "value": "Monitors high-privilege TCC grants"
    },
    "xpc_services_privileged": {
      "query": "SELECT la.label, la.program, la.path, la.run_at_load, la.keep_alive FROM launchd la WHERE la.path LIKE '%LaunchDaemons%' AND la.label NOT LIKE 'com.apple.%' ORDER BY la.label;",
      "interval": 3600,
      "description": "T1559: Non-Apple privileged XPC/LaunchDaemon services",
      "value": "Identifies third-party services running as root"
    }
  }
}
