# detections/sigma/persistence_launch_agent.yml
title: macOS LaunchAgent/LaunchDaemon Persistence
id: ct-persist-001
status: experimental
description: Detects creation of LaunchAgent or LaunchDaemon plist files for persistence
references:
  - https://attack.mitre.org/techniques/T1543/001/
  - https://attack.mitre.org/techniques/T1543/004/
author: The Eden's Sins Purple Team
date: 2025/01/15
tags:
  - attack.persistence
  - attack.t1543.001
  - attack.t1543.004
logsource:
  category: file_event
  product: macos
detection:
  selection_path:
    TargetFilename|contains:
      - '/Library/LaunchAgents/'
      - '/Library/LaunchDaemons/'
      - 'LaunchAgents/'
    TargetFilename|endswith: '.plist'
  filter_apple:
    TargetFilename|contains: 'com.apple.'
  filter_known:
    Image|contains:
      - '/usr/libexec/PlistBuddy'
      - '/System/Library/'
  condition: selection_path and not (filter_apple or filter_known)
level: high
falsepositives:
  - Legitimate software installation
  - MDM-managed deployments
---
title: macOS Dylib Hijacking Detection
id: ct-persist-002
status: experimental
description: Detects potential dylib hijacking via suspicious dylib loads or DYLD environment variables
references:
  - https://attack.mitre.org/techniques/T1574/004/
  - https://attack.mitre.org/techniques/T1574/006/
author: The Eden's Sins Purple Team
date: 2025/01/15
tags:
  - attack.persistence
  - attack.defense_evasion
  - attack.t1574.004
  - attack.t1574.006
logsource:
  category: process_creation
  product: macos
detection:
  selection_dyld_env:
    CommandLine|contains: 'DYLD_INSERT_LIBRARIES'
  selection_dyld_print:
    CommandLine|contains:
      - 'DYLD_PRINT_SEGMENTS'
      - 'DYLD_PRINT_BINDINGS'
      - 'DYLD_LIBRARY_PATH'
  selection_suspicious_dylib:
    ImageLoaded|contains:
      - '/tmp/'
      - '/var/tmp/'
      - '/Users/Shared/'
      - '/private/tmp/'
    ImageLoaded|endswith: '.dylib'
  condition: selection_dyld_env or selection_dyld_print or selection_suspicious_dylib
level: high
falsepositives:
  - Developer debugging
  - Legitimate code injection frameworks (Frida for testing)
